# rtsp-proxy.service
[Unit]
Description=RTSP Proxy Daemon
After=network.target

[Service]
# Most Rust apps run as simple foreground processes
Type=simple

ExecStart=/usr/local/bin/rtsp-proxy
# Optional: set if your app needs a working directory for config/logs
# WorkingDirectory=/opt/rtsp-proxy

# Restart policy
Restart=on-failure
RestartSec=5

# Run as non-root user (create the user first if you haven't)
User=rtspproxy
Group=rtspproxy

StandardOutput=journal
StandardError=journal

# Security hardening (optional but recommended)
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
